# deploy/redmine-max-notifier.service
; systemd-unit сервиса redmine-max-notifier.
;
; Установка:
;   sudo cp deploy/redmine-max-notifier.service /etc/systemd/system/
;   sudo systemctl daemon-reload
;   sudo systemctl enable --now redmine-max-notifier
;
; Подробности и подготовка окружения — README, раздел "Деплой".
;
; Комментарии в unit-файлах — английские директивы + русские пояснения:
; systemd читает только строки без ";" в начале.

[Unit]
Description=Redmine -> MAX notifier
Documentation=https://github.com/1eonardo9/redmine-max-notifier

; network-online, а не network: сервису нужен не поднятый интерфейс,
; а работающая сеть — на старте он резолвит адрес Redmine и собирает
; TLS-контекст для MAX.
After=network-online.target
Wants=network-online.target

; Redmine живёт на этой же машине, но в Docker (образ redmine:6, порт 80
; хоста проброшен в 3000 контейнера). Отдельного redmine.service не
; существует — контейнеры поднимает docker.service, поэтому цепляемся
; за него.
;
; Строка убирает пачку бессмысленных ошибок в логе при перезагрузке, но
; НЕ гарантирует готовность Redmine: docker-демон стартует за секунды,
; а Rails внутри контейнера поднимается ещё десятки секунд. Это и не
; нужно — цикл поллинга ловит исключения и повторит через минуту.
After=docker.service

[Service]
Type=simple

; Отдельный системный пользователь, без шелла и без прав на свой код.
; Создать: sudo useradd --system --no-create-home --shell /usr/sbin/nologin redmine-notifier
User=redmine-notifier
Group=redmine-notifier

; Код. Владелец — root, сервис его только читает: скомпрометированный
; сервис не должен уметь переписать собственные исходники.
WorkingDirectory=/opt/redmine-max-notifier

; Секреты вне рабочего дерева: .env рядом с кодом переживает git pull
; ровно до первого конфликта, а в бэкап каталога уезжает вместе с ним.
; Права: sudo chown root:redmine-notifier /etc/redmine-max-notifier/env
;        sudo chmod 640 /etc/redmine-max-notifier/env
EnvironmentFile=/etc/redmine-max-notifier/env

; StateDirectory создаст /var/lib/redmine-max-notifier с нужным владельцем
; и сам добавит его в ReadWritePaths. Именно сюда должен указывать
; DATABASE_URL — под ProtectSystem=strict больше писать некуда.
StateDirectory=redmine-max-notifier

ExecStart=/opt/redmine-max-notifier/.venv/bin/uvicorn \
    --factory redmine_max_notifier.web.app:create_app \
    --host 127.0.0.1 \
    --port 8000

; Сеть моргнула на старте, Redmine ещё не поднялся, битый путь к
; CA-bundle поправили — сервис вернётся сам.
Restart=on-failure
RestartSec=10s

; SIGINT, а не SIGTERM: uvicorn по SIGINT делает graceful shutdown и
; отрабатывает lifespan (гасит scheduler, закрывает клиентов и БД).
KillSignal=SIGINT
TimeoutStopSec=30

; Логи уезжают в journald: journalctl -u redmine-max-notifier -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redmine-max-notifier

; ── Ограничение прав ────────────────────────────────────────────────
; Сервису нужны ровно две вещи: исходящая сеть и своя БД.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
; AF_UNIX нужен не нашему коду, а glibc: резолвинг имён ходит в
; systemd-resolved (или nscd) через unix-сокет. Без него DNS до
; platform-api2.max.ru держится на фолбэке в /etc/resolv.conf — сработает
; или нет, зависит от хоста. Дешевле разрешить, чем ловить плавающее
; "иногда не резолвится".
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
